# cargos/cargos-release/files/etc/profile.d/pkgsrc.sh
# @PKGNAME@
# @EGDIR@/etc/profile.d/pkgsrc.sh

CVS_RSH=ssh
export CVS_RSH
