# cargos/cargos-release/files/etc/profile.d/sendmail.sh
# @PKGNAME@
# @EGDIR@/etc/profile.d/sendmail.sh

# SMTP relay host for busybox(1) sendmail (default: 127.0.0.1)
#SMTPHOST=

export SMTPHOST
